# validation/supervisor/chatterbox-flash-download.conf
[program:chatterbox-flash-download]
command=/workspace/Hybrid-tts/validation/run_chatterbox_flash_download.sh
autostart=false
autorestart=false
startsecs=0
stdout_logfile=/workspace/nano-flash-artifacts/g1/chatterbox-flash-download-supervisor.log
redirect_stderr=true
stdout_logfile_maxbytes=0
