# validation/supervisor/chatterbox-flash-setup.conf
[program:chatterbox-flash-setup]
command=/workspace/Hybrid-tts/validation/setup_chatterbox_flash.sh
autostart=false
autorestart=false
startsecs=0
stdout_logfile=/workspace/nano-flash-artifacts/g1/chatterbox-flash-setup-supervisor.log
redirect_stderr=true
stdout_logfile_maxbytes=0
