# validation/supervisor/chatterbox-nano-setup.conf
[program:chatterbox-nano-setup]
command=/workspace/Hybrid-tts/validation/setup_chatterbox_nano.sh
autostart=false
autorestart=false
startsecs=0
stdout_logfile=/workspace/nano-flash-artifacts/g1/chatterbox-nano-setup-supervisor.log
redirect_stderr=true
stdout_logfile_maxbytes=0
