# validation/supervisor/chatterbox-nano-profile.conf
[program:chatterbox-nano-profile]
command=/workspace/Hybrid-tts/validation/run_chatterbox_nano_profile.sh
autostart=false
autorestart=false
startsecs=0
stdout_logfile=/workspace/nano-flash-artifacts/g1/chatterbox-nano-profile-supervisor.log
redirect_stderr=true
stdout_logfile_maxbytes=0
