# validation/supervisor/final-adversarial.conf
[program:final-adversarial]
command=/bin/bash /workspace/Hybrid-tts/validation/run_final_adversarial.sh
directory=/workspace/Hybrid-tts
autostart=false
autorestart=false
startsecs=0
stdout_logfile=/workspace/nano-flash-artifacts/g8/final-adversarial-supervisor.log
stderr_logfile=/workspace/nano-flash-artifacts/g8/final-adversarial-supervisor.log
redirect_stderr=true
