# validation/supervisor/flashinfer-benchmark.conf
[program:flashinfer-benchmark]
command=/workspace/Hybrid-tts/validation/run_flashinfer_benchmark.sh
autostart=false
autorestart=false
startsecs=0
stdout_logfile=/workspace/nano-flash-artifacts/g7/flashinfer-benchmark-supervisor.log
redirect_stderr=true
stdout_logfile_maxbytes=0
