# validation/supervisor/flashinfer-setup.conf
[program:flashinfer-setup]
command=/workspace/Hybrid-tts/validation/setup_flashinfer.sh
autostart=false
autorestart=false
startsecs=0
stdout_logfile=/workspace/nano-flash-artifacts/g7/flashinfer-setup-supervisor.log
redirect_stderr=true
stdout_logfile_maxbytes=0
