# validation/supervisor/flash-quality.conf
[program:flash-quality]
command=/workspace/Hybrid-tts/validation/run_flash_quality_compare.sh
autostart=false
autorestart=false
startsecs=0
stdout_logfile=/workspace/nano-flash-artifacts/g7/flash-quality-supervisor.log
redirect_stderr=true
stdout_logfile_maxbytes=0
