# validation/supervisor/hybrid-interface.conf
[program:hybrid-interface]
command=/bin/bash /workspace/Hybrid-tts/validation/run_hybrid_interface_audit.sh
directory=/workspace/Hybrid-tts
autostart=false
autorestart=false
startsecs=0
stdout_logfile=/workspace/nano-flash-artifacts/g8/hybrid-interface-supervisor.log
stderr_logfile=/workspace/nano-flash-artifacts/g8/hybrid-interface-supervisor.log
redirect_stderr=true
