# validation/supervisor/librispeech-subset.conf
[program:librispeech-subset]
command=/workspace/Hybrid-tts/validation/run_librispeech_subset.sh
autostart=false
autorestart=false
startsecs=0
stdout_logfile=/workspace/nano-flash-artifacts/g4/librispeech-subset-supervisor.log
redirect_stderr=true
stdout_logfile_maxbytes=0
