# validation/supervisor/librispeech-encode.conf
[program:librispeech-encode]
command=/workspace/Hybrid-tts/validation/run_librispeech_encode.sh
autostart=false
autorestart=false
startsecs=0
stdout_logfile=/workspace/nano-flash-artifacts/g4/librispeech-encode-supervisor.log
redirect_stderr=true
stdout_logfile_maxbytes=0
