# validation/supervisor/moss-checkpoint-repro.conf
[program:moss-checkpoint-repro]
command=/workspace/Hybrid-tts/validation/run_moss_checkpoint_repro.sh
autostart=false
autorestart=false
startsecs=0
stdout_logfile=/workspace/nano-flash-artifacts/g3/moss-checkpoint-repro-supervisor.log
redirect_stderr=true
stdout_logfile_maxbytes=0
