# validation/supervisor/moss-lean-adaptation.conf
[program:moss-lean-adaptation]
command=/workspace/Hybrid-tts/validation/run_moss_lean_adaptation.sh
autostart=false
autorestart=false
startsecs=0
stdout_logfile=/workspace/nano-flash-artifacts/g4/moss-lean-adaptation-supervisor.log
redirect_stderr=true
stdout_logfile_maxbytes=0
