# validation/supervisor/moss-intelligibility.conf
[program:moss-intelligibility]
command=/workspace/Hybrid-tts/validation/run_moss_intelligibility.sh
autostart=false
autorestart=false
startsecs=0
stdout_logfile=/workspace/nano-flash-artifacts/g4/moss-intelligibility-supervisor.log
redirect_stderr=true
stdout_logfile_maxbytes=0
