# validation/supervisor/nano-flash-probe.conf
[program:nano-flash-probe]
command=/workspace/Hybrid-tts/validation/g0_supervisor_probe.sh
environment=NANO_FLASH_ARTIFACT_DIR="/workspace/nano-flash-artifacts/g0"
autostart=false
autorestart=false
startsecs=0
stdout_logfile=/dev/stdout
redirect_stderr=true
stdout_logfile_maxbytes=0
