# validation/supervisor/nano-renderer-cache.conf
[program:nano-renderer-cache]
command=/workspace/Hybrid-tts/validation/run_nano_renderer_cache.sh
autostart=false
autorestart=false
startsecs=0
stdout_logfile=/workspace/nano-flash-artifacts/g6/nano-renderer-cache-supervisor.log
redirect_stderr=true
stdout_logfile_maxbytes=0
