# validation/supervisor/nano-renderer-steps.conf
[program:nano-renderer-steps]
command=/workspace/Hybrid-tts/validation/run_nano_renderer_steps.sh
autostart=false
autorestart=false
startsecs=0
stdout_logfile=/workspace/nano-flash-artifacts/g6/nano-renderer-steps-supervisor.log
redirect_stderr=true
stdout_logfile_maxbytes=0
