# validation/supervisor/streaming-stress.conf
[program:streaming-stress]
command=/workspace/Hybrid-tts/validation/run_streaming_stress.sh
autostart=false
autorestart=false
startsecs=0
stdout_logfile=/workspace/nano-flash-artifacts/g5/streaming-stress-supervisor.log
redirect_stderr=true
stdout_logfile_maxbytes=0
